# vendor/osquery/packs/kubric-edr.conf
{
  "queries": {
    "listening_ports": {
      "query": "SELECT pid, port, protocol, address, path FROM listening_ports JOIN processes USING (pid);",
      "interval": 60,
      "description": "Tracks all listening network ports and their associated processes"
    },
    "unexpected_network_connections": {
      "query": "SELECT pid, remote_address, remote_port, local_port, protocol FROM process_open_sockets WHERE remote_address NOT IN ('0.0.0.0','::','127.0.0.1','::1') AND remote_port NOT IN (80,443,53);",
      "interval": 30,
      "description": "Detects unexpected outbound network connections"
    },
    "suspicious_processes": {
      "query": "SELECT p.pid, p.name, p.path, p.cmdline, p.uid, p.parent, pp.name AS parent_name FROM processes p LEFT JOIN processes pp ON p.parent = pp.pid WHERE p.path LIKE '%/tmp/%' OR p.path LIKE '%/dev/shm/%' OR p.cmdline LIKE '%-enc%' OR p.cmdline LIKE '%base64%';",
      "interval": 30,
      "description": "Finds processes running from temp locations or using suspicious command patterns"
    },
    "startup_items": {
      "query": "SELECT name, path, args, type FROM startup_items;",
      "interval": 300,
      "description": "Enumerates all startup items for persistence detection"
    },
    "suid_binaries": {
      "query": "SELECT path, username, groupname, permissions FROM file WHERE path LIKE '/usr/%' AND (permissions LIKE '%s%' OR permissions LIKE '%S%');",
      "interval": 3600,
      "description": "Lists SUID/SGID binaries that could be used for privilege escalation"
    },
    "crontab_entries": {
      "query": "SELECT command, path, minute, hour, day_of_month, month, day_of_week FROM crontab;",
      "interval": 300,
      "description": "Monitors crontab for unauthorized persistence entries"
    },
    "installed_software": {
      "query": "SELECT name, version, install_location FROM programs;",
      "interval": 3600,
      "description": "Tracks installed software inventory for change detection"
    },
    "kernel_modules": {
      "query": "SELECT name, used_by, status, address FROM kernel_modules;",
      "interval": 300,
      "description": "Monitors loaded kernel modules for rootkit detection"
    },
    "user_accounts": {
      "query": "SELECT uid, gid, username, description, directory, shell FROM users WHERE shell NOT IN ('/sbin/nologin','/bin/false','/usr/sbin/nologin');",
      "interval": 3600,
      "description": "Tracks user accounts with valid login shells"
    },
    "open_files": {
      "query": "SELECT p.name, p.pid, pof.fd, pof.path FROM process_open_files pof JOIN processes p ON pof.pid = p.pid WHERE pof.path LIKE '%.py' OR pof.path LIKE '%.sh' OR pof.path LIKE '/dev/%';",
      "interval": 60,
      "description": "Detects suspicious open file handles"
    }
  }
}
